# debian/postrm
#!/bin/bash
set -e

if [ "${HOST_BUILD_STAGE:-0}" -eq 1 ]; then
    exit 0
fi

if [ "$1" = "remove" ]; then
    dkms remove -m ethercat -v 1.7.1 --all >/dev/null 2>&1 || true
fi

if [ "$1" = "purge" ] || [ "$1" = "remove" ]; then
    PROFILE_FILE="/etc/profile.d/ethercat.sh"
    if [ -f "$PROFILE_FILE" ]; then
        rm -f "$PROFILE_FILE"
    fi
    if command -v systemctl >/dev/null 2>&1 && [ -d "/run/systemd/system" ]; then
        systemctl disable --now ethercat.service >/dev/null 2>&1 || true
        systemctl disable --now ethercat-first-boot.service >/dev/null 2>&1 || true
        systemctl daemon-reload || true
    fi
fi

if [ "$1" = "purge" ]; then
    rm -f /var/lib/ethercat/initialized
    rmdir /var/lib/ethercat 2>/dev/null || true
fi

if command -v ldconfig >/dev/null 2>&1; then
    ldconfig
fi

exit 0
