# debian/postinst
#!/bin/bash
set -e

if [ "${HOST_BUILD_STAGE:-0}" -eq 1 ]; then
    exit 0
fi

# DKMS: add, build, install kernel modules
NAME=ethercat
VERSION=1.7.1
if [ "$1" = "configure" ]; then
    if [ ! -d "/lib/modules/$(uname -r)/build" ]; then
        echo "ERROR: Kernel headers not found at /lib/modules/$(uname -r)/build" >&2
        echo "" >&2
        echo "DKMS needs kernel headers to compile the $NAME kernel modules." >&2
        echo "Install the matching kernel headers package, for example:" >&2
        echo "  Ubuntu/Debian:  apt install linux-headers-$(uname -r)" >&2
        echo "  Rockchip:       apt install linux-headers-current-rockchip-rk3588" >&2
        echo "  Or search:      apt search linux-headers | grep $(uname -r)" >&2
        exit 1
    fi
    dkms add -m $NAME -v $VERSION
    dkms build -m $NAME -v $VERSION
    dkms install -m $NAME -v $VERSION
fi

# PATH
PROFILE_DIR="/etc/profile.d"
PROFILE_FILE="$PROFILE_DIR/ethercat.sh"
if [ ! -d "$PROFILE_DIR" ]; then
    mkdir -p "$PROFILE_DIR"
fi
cat > "$PROFILE_FILE" <<'EOF'
# Add EtherCAT tools to PATH
if [ -d "/opt/ethercat/bin" ]; then
    export PATH="/opt/ethercat/bin:$PATH"
fi
EOF
chmod 644 "$PROFILE_FILE"

if command -v ldconfig >/dev/null 2>&1; then
    ldconfig
fi

# Auto-configure device module: igb for Rockchip kernels, generic for others
if [ -f /etc/ethercat.conf ]; then
    KERNEL_VER=$(uname -r)
    if echo "$KERNEL_VER" | grep -qi "rockchip\|rk35"; then
        sed -i 's/^DEVICE_MODULES=""/DEVICE_MODULES="igb"/' /etc/ethercat.conf
    else
        sed -i 's/^DEVICE_MODULES=""/DEVICE_MODULES="generic"/' /etc/ethercat.conf
    fi
fi

if command -v systemctl >/dev/null 2>&1 && [ -d "/run/systemd/system" ]; then
    mkdir -p /var/lib/ethercat
    systemctl daemon-reload || true
    systemctl enable ethercat.service >/dev/null 2>&1 || true
    systemctl enable ethercat-first-boot.service >/dev/null 2>&1 || true
fi

exit 0